# source/vendors/Ralink/RT6855A/inittab
::sysinit:/etc_ro/rcS
ttyS0::respawn:/bin/sh
